# cron.d/gotosleep
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# m h dom mon dow user	command
*  *  * * * root  bash -c 'DISPLAY=:0 /usr/bin/python3 /root/smartlock/smartlock.py &>> /smartlock.log'
#
